# systemd/livetrade-signal-monitor.timer
[Unit]
Description=Start Livetrade Signal Monitor Mon-Fri 09:15

[Timer]
OnCalendar=Mon..Fri 09:15
Persistent=true

[Install]
WantedBy=timers.target
